# systemd/.config/systemd/user/elephant-rescan.path
[Unit]
Description=Watch desktop-entry dirs and rescan elephant on change

[Path]
# Native apps + user/system flatpak exports. PathModified fires when a
# .desktop is added/removed (dir mtime changes).
PathModified=%h/.local/share/applications
PathModified=%h/.local/share/flatpak/exports/share/applications
PathModified=/var/lib/flatpak/exports/share/applications
Unit=elephant-rescan.service

[Install]
WantedBy=default.target
