# systemd/.config/systemd/user/rclone-vault-sync.path
[Unit]
Description=Sync vault when passwords.kdbx changes

[Path]
PathChanged=%h/vault/Passwords.kdbx

[Install]
WantedBy=default.target
