# systemd/.config/systemd/user/rclone-vault-sync.service
[Unit]
Description=Sync KeePassXC vault to Google Drive (rclone bisync)
After=network-online.target
Wants=network-online.target
# Fire a desktop notification if a sync run fails (so it can't break silently)
OnFailure=notify-failure@%n.service

[Service]
Type=oneshot
# First run ever (or after a break) must be done manually with --resync:
#   rclone bisync %h/vault gdrive:vault --resync
# After that, this keeps the two sides in sync.
#
# --force is REQUIRED here: the vault is a single file, so any KeePassXC save
# looks like "100% of files changed" and trips bisync's all-files-changed safety
# abort. --force bypasses that guard — which on its own is risky for a password
# DB (a corrupt/empty side would be propagated). The flags below make it safe:
#   --backup-dir1/2   before any file is deleted or overwritten on a side, the
#                     OLD version is moved to vault-backup/ on that same side, so
#                     a bad sync is always recoverable (last-known-good kept).
#   --conflict-resolve newer + --conflict-loser num  if BOTH sides changed
#                     (edited on two machines), keep the newer as the winner and
#                     preserve the older as Passwords.kdbx.conflict1 — never a
#                     silent loss of edits.
#   --recover / --resilient  auto-recover from interrupted/less-serious errors.
#   --check-access    RCLONE_TEST sentinel files exist on both sides; if a side
#                     reads empty (unmounted/auth-expired remote), bisync aborts
#                     instead of treating it as "delete everything". Re-create
#                     with: touch ~/vault/RCLONE_TEST && rclone copyto \
#                     ~/vault/RCLONE_TEST gdrive:vault/RCLONE_TEST  (+ a --resync).
ExecStart=/usr/bin/rclone bisync %h/vault gdrive:vault \
  --resilient --recover --force --check-access \
  --conflict-resolve newer --conflict-loser num --conflict-suffix conflict \
  --backup-dir1 %h/vault-backup --backup-dir2 gdrive:vault-backup \
  --create-empty-src-dirs
