# systemd/.config/systemd/user/rclone-vault-sync.service
[Unit]
Description=Sync KeePassXC vault to Google Drive (rclone bisync)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# First run ever (or after a break) must be done manually with --resync:
#   rclone bisync %h/vault gdrive:vault --resync
# After that, this keeps the two sides in sync. --resilient auto-recovers from
# interrupted runs. --force is REQUIRED here: the vault is a single file, so any
# KeePassXC save looks like "100% of files changed" and trips bisync's
# all-files-changed safety abort; --force lets the legitimate change through.
ExecStart=/usr/bin/rclone bisync %h/vault gdrive:vault --resilient --force --create-empty-src-dirs
