# systemd/.config/systemd/user/rclone-vault-sync.timer
[Unit]
Description=Periodic KeePassXC vault sync to Google Drive

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
Persistent=true

[Install]
WantedBy=timers.target
